# rpmbuild/SOURCES/fetch_blocked_ips.service
[Unit]
Description=Fetch and Upload Blocked IPs to GitHub
After=network.target

[Service]
ExecStart=/usr/local/bin/fetch_blocked_ips
Restart=on-failure

[Install]
WantedBy=multi-user.target
