# Fedora-RHEL-CentOS-AlmaLinux-RockyLinux/rpmbuild/SOURCES/honeypot-blocklist.service
[Unit]
Description=Honeypot Blocklist Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fetch_blocked_ips
ExecStartPost=/usr/local/bin/sync_blocklist
ExecStartPost=/usr/local/bin/update_blocklist

[Install]
WantedBy=multi-user.target
