# Fedora-RHEL-CentOS-AlmaLinux-RockyLinux/rpmbuild/SOURCES/sync_blocklist.service
[Unit]
Description=Sync Blocked IPs from GitHub and Apply to Firewalld
After=network.target

[Service]
ExecStart=/usr/local/bin/sync_blocklist
Restart=on-failure

[Install]
WantedBy=multi-user.target
