# Fedora-RHEL-CentOS-AlmaLinux-RockyLinux/etc/systemd/system/honeypot-blocklist.timer
[Unit]
Description=Run Honeypot Blocklist Service every hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
